# docker/supervisor.conf
[program:iogtmolo]
command = gunicorn --bind 0.0.0.0:8000 iogt.wsgi
environment = DJANGO_SETTINGS_MODULE="iogt.settings.docker"
directory = /deploy/
redirect_stderr = true
